# mcp/mcp.service.example
[Unit]
Description=BeeCount MCP Server
After=network.target

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/mcp
Environment="PATH=/path/to/mcp/mcp_env/bin"
ExecStart=/path/to/mcp/mcp_env/bin/python /path/to/mcp/run_server.py
Restart=always
RestartSec=10

# 日志
StandardOutput=journal
StandardError=journal
SyslogIdentifier=beecount-mcp

# 资源限制（根据服务器配置调整）
MemoryLimit=2G
CPUQuota=200%

[Install]
WantedBy=multi-user.target
